# glennergy-algoritm.service
[Unit]
Description=Glennergy Algoritm

[Service]
Type=oneshot
ExecStart=/usr/local/bin/Glennergy-Algoritm