# glennergy-spotpris.service
[Unit]
Description=Glennergy Spotpris

[Service]
Type=oneshot
ExecStart=/usr/local/bin/Glennergy-Spotpris