# glennergy-spotpris.timer
[Unit]
Description=Glennergy Spotpris Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=glennergy-spotpris.service

[Install]
WantedBy=timers.target