# glennergy-meteo.service
[Unit]
Description=Glennergy Meteo

[Service]
Type=oneshot
ExecStart=/usr/local/bin/Glennergy-Meteo
